# .config/systemd/user/clash_subscribe.service
[Unit]
Description=Update the clash config

[Service]
ExecStart=%h/bin/clash-subscribe.sh -w