# .config/systemd/user/clash_subscribe.timer
[Unit]
Description=Timer for update the clash config

[Timer]
OnCalendar=*-*-* *:0/10
Persistent=true

[Install]
WantedBy=timers.target