# .config/systemd/user/clash_subscribe.service
[Unit]
Description=Update the clash config

[Service]
EnvironmentFile=%h/.config/clash/subscribe_url.txt
ExecStart=%h/bin/clash-subscribe.sh -w
